# rsnapshot.cron
*  *   * * *   echo "crontab is working" &>/proc/1/fd/1
0 */4  * * *   /usr/src/app/job.sh alpha
30 3   * * *   /usr/src/app/job.sh beta
0  3   * * 1   /usr/src/app/job.sh gamma
30 2   1 * *   /usr/src/app/job.sh delta
